# yoomin_render.service
[Unit]
Description=YooMin LED Renderer
After=network.target

[Service]
Type=simple
# Replace /path/to/venv with your virtualenv path on the Raspberry Pi.
# This service sources the venv and runs render.py from the project directory.
User=root
WorkingDirectory=/home/vaibh/Documents/Yoo_Min_Art
# Use the venv python binary directly (absolute path) instead of sourcing.
ExecStart=/home/vaibh/Documents/Yoo_Min_Art/.venv/bin/python /home/vaibh/Documents/Yoo_Min_Art/render.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target